# upstart.conf
description "paperplane"
author  "Andrew Lobos and Ben Thomas"

start on startup
stop on shutdown

respawn

exec /usr/local/bin/node /var/local/paperplane/server.js >> /var/log/paperplane.log 2>&1
